# data/initrd/s390/etc/inittab
#id:0:initdefault:
id:1:initdefault:
# check system on startup
# first script to be executed if not booting in emergency (-b) mode
#si:I:bootwait:/sbin/halt
si:I:bootwait:/bin/netsetup
#si:I:bootwait:/bin/sync
1:1:respawn:/bin/bash
#1:1:once:/sbin/halt
# end of /etc/inittab
